# crontab
*/5 * * * * node /app/dist/registration/reconciliation-runner.js
